# infra/ai-social-bot.service
[Unit]
Description=AI Social Content Generator Telegram Bot
Documentation=https://github.com/s3qn/AI-Social-Content-Generator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sean
Group=sean
WorkingDirectory=/home/sean/ai-projects/AI-Social-Content-Generator
EnvironmentFile=/home/sean/ai-projects/AI-Social-Content-Generator/.env
# claude CLI lives under ~/.local/bin (a symlink to ~/.local/share/claude/versions/X).
# systemd's default PATH excludes this, so we prepend it to PATH.
Environment="PATH=/home/sean/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/sean/ai-projects/AI-Social-Content-Generator/.venv/bin/python -m ai_social_content_generator.telegram_bot.bot

# Restart policy: bring the bot back if it crashes or exits cleanly
# (the /restart admin command relies on this — it calls sys.exit(0))
Restart=always
RestartSec=3

# Don't restart endlessly if the bot is broken (e.g., bad code on deploy)
StartLimitIntervalSec=300
StartLimitBurst=5

# Logging: send stdout/stderr to systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-social-bot

# Security hardening (mild — bot needs network + filesystem access)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
